# config/supervisor/ODR-encoders.conf
# ENERGY98
[program:odr-padencoder-63a4992d-b933-4136-90ef-34512340964e]
command=/usr/local/bin/odr-padenc
 --dir=/home/odr/ODR-mmbTools/config/mot/P01/slide
 --sleep=0
 --dls=/home/odr/ODR-mmbTools/config/mot/P01/INFO.dls
 --output=63a4992d-b933-4136-90ef-34512340964e
 --label=12
 --label-ins=1200
autostart=false
autorestart=true
priority=10
user=odr
group=odr
redirect_stderr=true
stdout_logfile=/var/log/supervisor/odr-padencoder-63a4992d-b933-4136-90ef-34512340964e.log

# ENERGY98
[program:odr-audioencoder-63a4992d-b933-4136-90ef-34512340964e]
command=/usr/local/bin/odr-audioenc
 --gst-uri=http://listen.livestreamingservice.com:8800
 --drift-comp
 --silence=60
 --bitrate=88
 --rate=48000
 --channels=2
 --ps
 --pad=34
 --pad-socket=63a4992d-b933-4136-90ef-34512340964e
 --write-icy-text=/home/odr/ODR-mmbTools/config/mot/P01/INFO.dls
 --edi=tcp://127.0.0.1:9001
 --stats=/var/tmp/P01.stats
autostart=false
autorestart=true
priority=10
user=odr
group=odr
redirect_stderr=true
stdout_logfile=/var/log/supervisor/odr-audioencoder-63a4992d-b933-4136-90ef-34512340964e.log

# POWER
[program:odr-padencoder-91362c6e-77a7-4aed-9527-02b695484e2f]
command=/usr/local/bin/odr-padenc
 --dir=/home/odr/ODR-mmbTools/config/mot/P02/slide
 --sleep=0
 --dls=/home/odr/ODR-mmbTools/config/mot/P02/INFO.dls
 --output=91362c6e-77a7-4aed-9527-02b695484e2f
 --label=12
 --label-ins=1200
autostart=false
autorestart=true
priority=10
user=odr
group=odr
redirect_stderr=true
stdout_logfile=/var/log/supervisor/odr-padencoder-91362c6e-77a7-4aed-9527-02b695484e2f.log

# POWER
[program:odr-audioencoder-91362c6e-77a7-4aed-9527-02b695484e2f]
command=/usr/local/bin/odr-audioenc
 --gst-uri=http://listen.livestreamingservice.com:8128
 --drift-comp
 --silence=60
 --bitrate=88
 --rate=48000
 --channels=2
 --ps
 --pad=34
 --pad-socket=91362c6e-77a7-4aed-9527-02b695484e2f
 --write-icy-text=/home/odr/ODR-mmbTools/config/mot/P02/INFO.dls
 --edi=tcp://127.0.0.1:9002
 --stats=/var/tmp/P02.stats
autostart=false
autorestart=true
priority=10
user=odr
group=odr
redirect_stderr=true
stdout_logfile=/var/log/supervisor/odr-audioencoder-91362c6e-77a7-4aed-9527-02b695484e2f.log

# UKTOP40
[program:odr-padencoder-8ad19122-6b76-40bf-a9d8-19f5435f04b4]
command=/usr/local/bin/odr-padenc
 --dir=/home/odr/ODR-mmbTools/config/mot/P03/slide
 --sleep=0
 --dls=/home/odr/ODR-mmbTools/config/mot/P03/INFO.dls
 --output=8ad19122-6b76-40bf-a9d8-19f5435f04b4
 --label=12
 --label-ins=1200
autostart=false
autorestart=true
priority=10
user=odr
group=odr
redirect_stderr=true
stdout_logfile=/var/log/supervisor/odr-padencoder-8ad19122-6b76-40bf-a9d8-19f5435f04b4.log

# UKTOP40
[program:odr-audioencoder-8ad19122-6b76-40bf-a9d8-19f5435f04b4]
command=/usr/local/bin/odr-audioenc
 --gst-uri=http://listen.livestreamingservice.com:8070
 --drift-comp
 --silence=60
 --bitrate=88
 --rate=48000
 --channels=2
 --ps
 --pad=34
 --pad-socket=8ad19122-6b76-40bf-a9d8-19f5435f04b4
 --write-icy-text=/home/odr/ODR-mmbTools/config/mot/P03/INFO.dls
 --edi=tcp://127.0.0.1:9003
 --stats=/var/tmp/P03.stats
autostart=false
autorestart=true
priority=10
user=odr
group=odr
redirect_stderr=true
stdout_logfile=/var/log/supervisor/odr-audioencoder-8ad19122-6b76-40bf-a9d8-19f5435f04b4.log

# ROCK40
[program:odr-padencoder-cbd9385b-8499-4c0b-96b7-f2734651214b]
command=/usr/local/bin/odr-padenc
 --dir=/home/odr/ODR-mmbTools/config/mot/P04/slide
 --sleep=0
 --dls=/home/odr/ODR-mmbTools/config/mot/P04/INFO.dls
 --output=cbd9385b-8499-4c0b-96b7-f2734651214b
 --label=12
 --label-ins=1200
autostart=false
autorestart=true
priority=10
user=odr
group=odr
redirect_stderr=true
stdout_logfile=/var/log/supervisor/odr-padencoder-cbd9385b-8499-4c0b-96b7-f2734651214b.log

# ROCK40
[program:odr-audioencoder-cbd9385b-8499-4c0b-96b7-f2734651214b]
command=/usr/local/bin/odr-audioenc
 --gst-uri=http://listen.livestreamingservice.com:8028
 --drift-comp
 --silence=60
 --bitrate=88
 --rate=48000
 --channels=2
 --ps
 --pad=34
 --pad-socket=cbd9385b-8499-4c0b-96b7-f2734651214b
 --write-icy-text=/home/odr/ODR-mmbTools/config/mot/P04/INFO.dls
 --edi=tcp://127.0.0.1:9004
 --stats=/var/tmp/P04.stats
autostart=false
autorestart=true
priority=10
user=odr
group=odr
redirect_stderr=true
stdout_logfile=/var/log/supervisor/odr-audioencoder-cbd9385b-8499-4c0b-96b7-f2734651214b.log

# CLASSIC HITS
[program:odr-padencoder-3b96520b-a58e-410a-a5d3-b4134445a458]
command=/usr/local/bin/odr-padenc
 --dir=/home/odr/ODR-mmbTools/config/mot/P05/slide
 --sleep=0
 --dls=/home/odr/ODR-mmbTools/config/mot/P05/INFO.dls
 --output=3b96520b-a58e-410a-a5d3-b4134445a458
 --label=12
 --label-ins=1200
autostart=false
autorestart=true
priority=10
user=odr
group=odr
redirect_stderr=true
stdout_logfile=/var/log/supervisor/odr-padencoder-3b96520b-a58e-410a-a5d3-b4134445a458.log

# CLASSIC HITS
[program:odr-audioencoder-3b96520b-a58e-410a-a5d3-b4134445a458]
command=/usr/local/bin/odr-audioenc
 --gst-uri=http://listen.livestreamingservice.com:8132
 --drift-comp
 --silence=60
 --bitrate=88
 --rate=48000
 --channels=2
 --ps
 --pad=34
 --pad-socket=3b96520b-a58e-410a-a5d3-b4134445a458
 --write-icy-text=/home/odr/ODR-mmbTools/config/mot/P05/INFO.dls
 --edi=tcp://127.0.0.1:9005
 --stats=/var/tmp/P05.stats
autostart=false
autorestart=true
priority=10
user=odr
group=odr
redirect_stderr=true
stdout_logfile=/var/log/supervisor/odr-audioencoder-3b96520b-a58e-410a-a5d3-b4134445a458.log

# AWESOME80
[program:odr-padencoder-0d741242-79ec-4ae8-82b5-0bdda3deb6f2]
command=/usr/local/bin/odr-padenc
 --dir=/home/odr/ODR-mmbTools/config/mot/P06/slide
 --sleep=0
 --dls=/home/odr/ODR-mmbTools/config/mot/P06/INFO.dls
 --output=0d741242-79ec-4ae8-82b5-0bdda3deb6f2
 --label=12
 --label-ins=1200
autostart=false
autorestart=true
priority=10
user=odr
group=odr
redirect_stderr=true
stdout_logfile=/var/log/supervisor/odr-padencoder-0d741242-79ec-4ae8-82b5-0bdda3deb6f2.log

# AWESOME80
[program:odr-audioencoder-0d741242-79ec-4ae8-82b5-0bdda3deb6f2]
command=/usr/local/bin/odr-audioenc
 --gst-uri=http://listen.livestreamingservice.com:8000
 --drift-comp
 --silence=60
 --bitrate=88
 --rate=48000
 --channels=2
 --ps
 --pad=34
 --pad-socket=0d741242-79ec-4ae8-82b5-0bdda3deb6f2
 --write-icy-text=/home/odr/ODR-mmbTools/config/mot/P06/INFO.dls
 --edi=tcp://127.0.0.1:9006
 --stats=/var/tmp/P06.stats
autostart=false
autorestart=true
priority=10
user=odr
group=odr
redirect_stderr=true
stdout_logfile=/var/log/supervisor/odr-audioencoder-0d741242-79ec-4ae8-82b5-0bdda3deb6f2.log

# THE VIBE
[program:odr-padencoder-a0fa049b-a726-43eb-bf63-4f725e899d84]
command=/usr/local/bin/odr-padenc
 --dir=/home/odr/ODR-mmbTools/config/mot/P07/slide
 --sleep=0
 --dls=/home/odr/ODR-mmbTools/config/mot/P07/INFO.dls
 --output=a0fa049b-a726-43eb-bf63-4f725e899d84
 --label=12
 --label-ins=1200
autostart=false
autorestart=true
priority=10
user=odr
group=odr
redirect_stderr=true
stdout_logfile=/var/log/supervisor/odr-padencoder-a0fa049b-a726-43eb-bf63-4f725e899d84.log

# THE VIBE
[program:odr-audioencoder-a0fa049b-a726-43eb-bf63-4f725e899d84]
command=/usr/local/bin/odr-audioenc
 --gst-uri=http://listen.livestreamingservice.com:8074
 --drift-comp
 --silence=60
 --bitrate=88
 --rate=48000
 --channels=2
 --ps
 --pad=34
 --pad-socket=a0fa049b-a726-43eb-bf63-4f725e899d84
 --write-icy-text=/home/odr/ODR-mmbTools/config/mot/P07/INFO.dls
 --edi=tcp://127.0.0.1:9007
 --stats=/var/tmp/P07.stats
autostart=false
autorestart=true
priority=10
user=odr
group=odr
redirect_stderr=true
stdout_logfile=/var/log/supervisor/odr-audioencoder-a0fa049b-a726-43eb-bf63-4f725e899d84.log

# THE BEAT
[program:odr-padencoder-30d42a0a-d88a-4215-9ab6-207d63b77100]
command=/usr/local/bin/odr-padenc
 --dir=/home/odr/ODR-mmbTools/config/mot/P08/slide
 --sleep=0
 --dls=/home/odr/ODR-mmbTools/config/mot/P08/INFO.dls
 --output=30d42a0a-d88a-4215-9ab6-207d63b77100
 --label=12
 --label-ins=1200
autostart=false
autorestart=true
priority=10
user=odr
group=odr
redirect_stderr=true
stdout_logfile=/var/log/supervisor/odr-padencoder-30d42a0a-d88a-4215-9ab6-207d63b77100.log

# THE BEAT
[program:odr-audioencoder-30d42a0a-d88a-4215-9ab6-207d63b77100]
command=/usr/local/bin/odr-audioenc
 --gst-uri=http://listen.livestreamingservice.com:8054
 --drift-comp
 --silence=60
 --bitrate=88
 --rate=48000
 --channels=2
 --ps
 --pad=34
 --pad-socket=30d42a0a-d88a-4215-9ab6-207d63b77100
 --write-icy-text=/home/odr/ODR-mmbTools/config/mot/P08/INFO.dls
 --edi=tcp://127.0.0.1:9008
 --stats=/var/tmp/P08.stats
autostart=false
autorestart=true
priority=10
user=odr
group=odr
redirect_stderr=true
stdout_logfile=/var/log/supervisor/odr-audioencoder-30d42a0a-d88a-4215-9ab6-207d63b77100.log

# ENERGY93
[program:odr-padencoder-9b0eb1fe-8967-4c2a-8175-9e12002db557]
command=/usr/local/bin/odr-padenc
 --dir=/home/odr/ODR-mmbTools/config/mot/P09/slide
 --sleep=0
 --dls=/home/odr/ODR-mmbTools/config/mot/P09/INFO.dls
 --output=9b0eb1fe-8967-4c2a-8175-9e12002db557
 --label=12
 --label-ins=1200
autostart=false
autorestart=true
priority=10
user=odr
group=odr
redirect_stderr=true
stdout_logfile=/var/log/supervisor/odr-padencoder-9b0eb1fe-8967-4c2a-8175-9e12002db557.log

# ENERGY93
[program:odr-audioencoder-9b0eb1fe-8967-4c2a-8175-9e12002db557]
command=/usr/local/bin/odr-audioenc
 --gst-uri=http://listen.livestreamingservice.com:8044
 --drift-comp
 --silence=60
 --bitrate=88
 --rate=48000
 --channels=2
 --ps
 --pad=34
 --pad-socket=9b0eb1fe-8967-4c2a-8175-9e12002db557
 --write-icy-text=/home/odr/ODR-mmbTools/config/mot/P09/INFO.dls
 --edi=tcp://127.0.0.1:9009
 --stats=/var/tmp/P09.stats
autostart=false
autorestart=true
priority=10
user=odr
group=odr
redirect_stderr=true
stdout_logfile=/var/log/supervisor/odr-audioencoder-9b0eb1fe-8967-4c2a-8175-9e12002db557.log

# MELLOW GOLD
[program:odr-padencoder-86d4545b-50b6-4dab-8560-07354906b529]
command=/usr/local/bin/odr-padenc
 --dir=/home/odr/ODR-mmbTools/config/mot/P10/slide
 --sleep=0
 --dls=/home/odr/ODR-mmbTools/config/mot/P10/INFO.dls
 --output=86d4545b-50b6-4dab-8560-07354906b529
 --label=12
 --label-ins=1200
autostart=false
autorestart=true
priority=10
user=odr
group=odr
redirect_stderr=true
stdout_logfile=/var/log/supervisor/odr-padencoder-86d4545b-50b6-4dab-8560-07354906b529.log

# MELLOW GOLD
[program:odr-audioencoder-86d4545b-50b6-4dab-8560-07354906b529]
command=/usr/local/bin/odr-audioenc
 --gst-uri=http://listen.livestreamingservice.com:8060
 --drift-comp
 --silence=60
 --bitrate=88
 --rate=48000
 --channels=2
 --ps
 --pad=34
 --pad-socket=86d4545b-50b6-4dab-8560-07354906b529
 --write-icy-text=/home/odr/ODR-mmbTools/config/mot/P10/INFO.dls
 --edi=tcp://127.0.0.1:9010
 --stats=/var/tmp/P10.stats
autostart=false
autorestart=true
priority=10
user=odr
group=odr
redirect_stderr=true
stdout_logfile=/var/log/supervisor/odr-audioencoder-86d4545b-50b6-4dab-8560-07354906b529.log

# STAR90
[program:odr-padencoder-f72fcb25-86db-486b-90d3-09dc96cf5a32]
command=/usr/local/bin/odr-padenc
 --dir=/home/odr/ODR-mmbTools/config/mot/P11/slide
 --sleep=0
 --dls=/home/odr/ODR-mmbTools/config/mot/P11/INFO.dls
 --output=f72fcb25-86db-486b-90d3-09dc96cf5a32
 --label=12
 --label-ins=1200
autostart=false
autorestart=true
priority=10
user=odr
group=odr
redirect_stderr=true
stdout_logfile=/var/log/supervisor/odr-padencoder-f72fcb25-86db-486b-90d3-09dc96cf5a32.log

# STAR90
[program:odr-audioencoder-f72fcb25-86db-486b-90d3-09dc96cf5a32]
command=/usr/local/bin/odr-audioenc
 --gst-uri=http://listen.livestreamingservice.com:8012
 --drift-comp
 --silence=60
 --bitrate=88
 --rate=48000
 --channels=2
 --ps
 --pad=34
 --pad-socket=f72fcb25-86db-486b-90d3-09dc96cf5a32
 --write-icy-text=/home/odr/ODR-mmbTools/config/mot/P11/INFO.dls
 --edi=tcp://127.0.0.1:9011
 --stats=/var/tmp/P11.stats
autostart=false
autorestart=true
priority=10
user=odr
group=odr
redirect_stderr=true
stdout_logfile=/var/log/supervisor/odr-audioencoder-f72fcb25-86db-486b-90d3-09dc96cf5a32.log

# THE BREEZE
[program:odr-padencoder-3a9e3453-f98d-4a40-9cc3-75f8a26e42b1]
command=/usr/local/bin/odr-padenc
 --dir=/home/odr/ODR-mmbTools/config/mot/P12/slide
 --sleep=0
 --dls=/home/odr/ODR-mmbTools/config/mot/P12/INFO.dls
 --output=3a9e3453-f98d-4a40-9cc3-75f8a26e42b1
 --label=12
 --label-ins=1200
autostart=false
autorestart=true
priority=10
user=odr
group=odr
redirect_stderr=true
stdout_logfile=/var/log/supervisor/odr-padencoder-3a9e3453-f98d-4a40-9cc3-75f8a26e42b1.log

# THE BREEZE
[program:odr-audioencoder-3a9e3453-f98d-4a40-9cc3-75f8a26e42b1]
command=/usr/local/bin/odr-audioenc
 --gst-uri=http://listen.livestreamingservice.com:8004
 --drift-comp
 --silence=60
 --bitrate=88
 --rate=48000
 --channels=2
 --ps
 --pad=34
 --pad-socket=3a9e3453-f98d-4a40-9cc3-75f8a26e42b1
 --write-icy-text=/home/odr/ODR-mmbTools/config/mot/P12/INFO.dls
 --edi=tcp://127.0.0.1:9012
 --stats=/var/tmp/P12.stats
autostart=false
autorestart=true
priority=10
user=odr
group=odr
redirect_stderr=true
stdout_logfile=/var/log/supervisor/odr-audioencoder-3a9e3453-f98d-4a40-9cc3-75f8a26e42b1.log

# CLASSICAL
[program:odr-padencoder-1b8d0058-3f43-46e1-848c-debd31836dc2]
command=/usr/local/bin/odr-padenc
 --dir=/home/odr/ODR-mmbTools/config/mot/P13/slide
 --sleep=0
 --dls=/home/odr/ODR-mmbTools/config/mot/P13/INFO.dls
 --output=1b8d0058-3f43-46e1-848c-debd31836dc2
 --label=12
 --label-ins=1200
autostart=false
autorestart=true
priority=10
user=odr
group=odr
redirect_stderr=true
stdout_logfile=/var/log/supervisor/odr-padencoder-1b8d0058-3f43-46e1-848c-debd31836dc2.log

# CLASSICAL
[program:odr-audioencoder-1b8d0058-3f43-46e1-848c-debd31836dc2]
command=/usr/local/bin/odr-audioenc
 --gst-uri=http://listen.livestreamingservice.com/181-classical_128k.mp3
 --drift-comp
 --silence=60
 --bitrate=88
 --rate=48000
 --channels=2
 --ps
 --pad=34
 --pad-socket=1b8d0058-3f43-46e1-848c-debd31836dc2
 --write-icy-text=/home/odr/ODR-mmbTools/config/mot/P13/INFO.dls
 --edi=tcp://127.0.0.1:9013
 --stats=/var/tmp/P13.stats
autostart=false
autorestart=true
priority=10
user=odr
group=odr
redirect_stderr=true
stdout_logfile=/var/log/supervisor/odr-audioencoder-1b8d0058-3f43-46e1-848c-debd31836dc2.log

